# 10_processes/supervisord.conf
[program:daemon]
command=/root/daemon.py ; the program (relative uses PATH, can take args)
autorestart=true ; retstart at unexpected quit (default: true)